# app.conf
[program:gunicorn]
command:gunicorn microblog:app
directory: /home/ubuntu
user: ubuntu
autostart:true

[program:microblog]
command:microblog.py runserver
directory: /home/ubuntu
autostart:false
